# docker/php/supervisor/queue.conf
[program:rabbit]
process_name=%(program_name)s_%(process_num)02d
command=php /usr/src/app/bin/console rabbitmq:consumer client_notify
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/usr/src/supervisor/worker.log
priority=2
startsecs=20

[program:php-fpm]
command=php-fpm -R -F -c /etc/php/7.4/fpm/php-fpm.conf
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=1
